# iqfanda-host-agent/rootfs/etc/services.d/iqf-host-agent/finish
#!/usr/bin/with-contenv bashio

bashio::log.info "TNG IQ FANDA Host Agent byl ukončen."
