# iqfanda-host-agent/rootfs/etc/services.d/iqf-host-api/finish
#!/usr/bin/with-contenv bashio

bashio::log.info "TNG IQ FANDA Host API bylo zastaveno."
